# deploy/systemd/borofone-staging.service
[Unit]
Description=Borofone Chat Staging
After=network.target

[Service]
User=root
WorkingDirectory=/opt/borofone-chat-staging
EnvironmentFile=/opt/borofone-chat-staging/.env
ExecStart=/opt/borofone-chat-staging/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
